# docker-state-saver/docker-state-saver.service
[Unit]
Description=Docker State Saver - Persist and Restore Running Containers
# This service must run after the Docker socket is available.
# It must run before the system shuts down.
Requires=docker.socket
After=docker.socket
Before=shutdown.target

[Service]
Type=oneshot
RemainAfterExit=yes

# On system start, run the 'restore' command.
ExecStart=/usr/local/bin/docker-state-saver.sh restore

# On system stop/shutdown, run the 'save' command.
ExecStop=/usr/local/bin/docker-state-saver.sh save

[Install]
# Enable this service for the standard multi-user system state.
WantedBy=multi-user.target
